# linux/jumpproxy.service
[Unit]
Description=JumpProxy XVPN SDK public SOCKS relay
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/JumpProxyLinuxWine
EnvironmentFile=-/etc/jumpproxy.env
ExecStart=/bin/bash /opt/JumpProxyLinuxWine/linux/run_xvpn_wine_proxy.sh --with-tun2socks --with-iphlpapi-shim
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
