# linux/jumpproxy.service
[Unit]
Description=JumpProxy XVPN SDK public SOCKS relay
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/JumpProxyLinuxWine
Environment=PUBLIC_PORT=10880
Environment=PUBLIC_LISTEN=0.0.0.0
Environment=MODE=ad
Environment=HEALTH_INTERVAL=120
Environment=HEALTH_FAILURES=5
Environment=HEALTH_PROBES=3
Environment=HEALTH_PROBE_MAX_FAILURES=1
Environment=RECONNECT_DELAY=10
Environment=START_WAIT=45
Environment=STABILITY_PROBES=2
Environment=STABILITY_MAX_FAILURES=1
Environment=PUBLIC_MAX_CONNECTIONS=16
Environment=PUBLIC_UPSTREAM_RETRIES=12
Environment=PUBLIC_STREAM_RETRIES=4
Environment=PUBLIC_CLIENT_FAILOVER_ATTEMPTS=4
Environment=PUBLIC_CLIENT_FAILOVER_WAIT_SECONDS=120
Environment=PUBLIC_CONNECT_FAILURES_BEFORE_ROTATE=1
Environment=SKIP_HEALTH_WHEN_PUBLIC_ACTIVE_SECONDS=300
Environment=DEFER_ROTATE_WHEN_PUBLIC_ACTIVE_SECONDS=300
Environment=MAX_ROTATE_DEFER_SECONDS=900
Environment=HEALTH_VIA_PUBLIC=0
Environment=XVPN_DISABLE_UDP=1
Environment=XVPN_WINE_DLL_OVERRIDES=iphlpapi,fwpuclnt=n,b
ExecStart=/bin/bash /opt/JumpProxyLinuxWine/linux/run_xvpn_wine_proxy.sh --with-tun2socks --with-iphlpapi-shim
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
